# thirdparty/60-als.rules
ACTION=="change", KERNEL=="ACPI0008:00", RUN+="/home/oktay/Projekte/lightndark/lightndark/update.py /home/oktay/Projekte/lightndark/lightndark/config.ini"
